# scripts/com.assistant.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- ラベル（一意な識別子） -->
    <key>Label</key>
    <string>com.assistant.daemon</string>

    <!-- 実行するプログラム -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/YOUR_USERNAME/ai-assistant/assistant-daemon.js</string>
    </array>

    <!-- 作業ディレクトリ -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/ai-assistant</string>

    <!-- 環境変数 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/YOUR_USERNAME</string>
    </dict>

    <!-- ログイン時に自動起動 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- クラッシュ時に自動再起動 -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- 標準出力ログ -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/ai-assistant/logs/daemon.log</string>

    <!-- エラー出力ログ -->
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/ai-assistant/logs/daemon-error.log</string>

    <!-- CPU 使用率制限（オプション） -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- 開始遅延（秒）- 起動直後の負荷軽減 -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
